# base/snail-files/files/etc/profile.d/ls.sh
#!/bin/sh

alias .='cd .'
alias ..='cd ..'
alias ...='cd ../../'

alias  l='ls -CF'
alias la='ls -A'
alias lf='ls --color=auto -alFh'
alias ll='ls --color=auto -alF'
alias ls='ls --color=auto'


