# base/snail-files/files/etc/profile.d/ps.sh
#!/bin/sh

alias  af='/bin/ps af'
alias ps1='/bin/ps faxo "%U %t %p %a" '
alias ps2='/bin/ps facux '
alias psg='/bin/busybox ps | grep -n'


