# base/snail-files/files/etc/profile.d/vim.sh
#!/bin/sh

[ -x /usr/bin/vim ] && \
{
	alias  vi='/usr/bin/vim'
	alias vim='/usr/bin/vim'
}

